# config/crontab
# PATH=/opt/bin
*/1 * * * * cd /opt/var/iyyov && /opt/jruby/gems/gems/iyyov-1.0.2-java/init/iyyov jobs.rb
